# recipes-core/udev/udev-extraconf/20-drm-fb1-auto-show-logo.rules
ACTION=="add", SUBSYSTEM=="graphics", KERNEL=="fb1", RUN+="/usr/bin/fbi -d /dev/fb1 -a -T 7 /usr/share/densitron/quixant.png"
ACTION=="change", SUBSYSTEM=="graphics", KERNEL=="fb1", RUN+="/usr/bin/fbi -d /dev/fb1 -a -T 7 /usr/share/densitron/quixant.png"
